# systemd/ayugram.service
[Unit]
Description=Watch AyuGram and kill it every 180 seconds
After=graphical.target

[Service]
Type=simple
ExecStart=/bin/bash -c "while true; do sleep 180; if pgrep -x AyuGram > /dev/null; then echo \"\$(date): Killing AyuGram\"; pkill -x AyuGram; fi; done"
Restart=always
